# launchd/com.vure.temu-report.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.vure.temu-report</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/vure/ReportDalily/scripts/run-temu-report-all.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/vure/ReportDalily</string>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>23</integer>
      <key>Minute</key>
      <integer>50</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/vure/ReportDalily/temu-reports/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/vure/ReportDalily/temu-reports/launchd.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
